# systemd/lawnberry-backend.service
[Unit]
Description=LawnBerry Pi v2 Backend API Service
Documentation=https://github.com/acredsfan/lawnberry_pi
After=network-online.target
Wants=network-online.target
Requires=lawnberry-database.service

[Service]
Type=exec
User=pi
Group=pi
WorkingDirectory=/apps/lawnberry-pi
Environment=PYTHONPATH=/apps/lawnberry-pi
Environment=LAWNBERRY_ENV=production
Environment=LAWNBERRY_CONFIG_DIR=/apps/lawnberry-pi/config
Environment=LAWNBERRY_LOG_DIR=/apps/lawnberry-pi/logs
Environment=LAWNBERRY_DATA_DIR=/apps/lawnberry-pi/data
Environment=SIM_MODE=0
Environment=PYTHONUNBUFFERED=1
# Auth disabled project-wide (LAN-only single-operator deployment, no
# internet exposure): login, operator-auth-gated writes, and the manual
# control unlock step all no-op. See api/deps.py::require_operator_auth
# and api/routers/auth.py's manual-unlock handlers.
Environment=OPERATOR_AUTH_REQUIRED=0
# Auto-upload recorded mowing sessions to the Thor training server.
# Set THOR_UPLOAD_ENABLED=1 and THOR_BASE_URL to Thor's receiver (host:port).
Environment=THOR_UPLOAD_ENABLED=0
Environment=THOR_BASE_URL=http://192.168.1.64:15454
# Load runtime secrets and configuration (e.g., NTRIP_*) from .env if present
EnvironmentFile=-/apps/lawnberry-pi/.env

# Main service command
ExecStart=/apps/lawnberry-pi/.venv/bin/uvicorn backend.src.main:app --host 0.0.0.0 --port 8081 --workers 1

# Health checks and monitoring
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Process management
Restart=always
RestartSec=10
RestartPreventExitStatus=2

# Security restrictions (ARM64/Pi OS specific)
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
# App root must be writable: lgpio creates its .lgd-nfy-* notification FIFO in the
# process CWD; with only subdirs writable, ultrasonic GPIO init fails (ENOENT).
ReadWritePaths=/apps/lawnberry-pi
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
UMask=007

# Hardware access (required for GPIO, I2C, SPI)
SupplementaryGroups=gpio i2c spi dialout

# Resource limits
LimitNOFILE=4096
MemoryMax=512M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-backend

[Install]
WantedBy=multi-user.target
